# deployment/etc/init/gunicorn_sudoku.conf
description "upstart configuration for gunicorn sudoku test project"

start on net-device-up
stop on shutdown

respawn
setuid ubuntu
setgid ubuntu

exec /var/www/sudoku/deployment/start_gunicorn.sh
